# infra/launchd/com.trader.weekly-enforcing-canary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  v3.73.26 — Weekly ENFORCING canary.

  Sundays at 15:00 UTC. Re-runs the synthetic -13% drawdown drill
  in-process and verifies that ENFORCING still mutates targets
  correctly. If the brake silently broke during the 30-run window,
  this fires Slack + email and logs a FAIL row in
  docs/ENFORCING_CANARY_LOG.md. The canary does NOT touch the real
  journal and CANNOT break the 30-run streak.

  Sleep-resilience: paired StartInterval (3 hours) so a missed
  Sunday window is caught the next time the laptop wakes. The
  canary is idempotent — over-firing just adds another green row
  to the log.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.trader.weekly-enforcing-canary</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd /Users/richardchen/trader &amp;&amp; set -a &amp;&amp; [ -f .env ] &amp;&amp; source .env; set +a; exec .venv/bin/python scripts/weekly_enforcing_canary.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <!-- Weekday 0 = Sunday on launchd. 15:00 UTC = ~10am ET. -->
    <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <!-- Sleep-resilience: 3-hour interval ensures missed Sunday fire
       is picked up on next wake. Canary is idempotent (just appends
       to the log). -->
  <key>StartInterval</key>
  <integer>10800</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/richardchen/Library/Logs/trader-weekly-enforcing-canary.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/richardchen/Library/Logs/trader-weekly-enforcing-canary.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>/Users/richardchen</string>
    <key>PATH</key><string>/Users/richardchen/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/richardchen/trader</string>

  <key>ProcessType</key>
  <string>Adaptive</string>
</dict>
</plist>
